# scripts/strategy-dashboard.service
[Unit]
Description=Strategy Console Dashboard
After=network.target

[Service]
Type=simple
User=robpitman1982
WorkingDirectory=/home/robpitman1982/python-master-strategy-creator
Environment=STRATEGY_CONSOLE_STORAGE=/home/robpitman1982/strategy_console_storage
ExecStart=/home/robpitman1982/python-master-strategy-creator/venv/bin/streamlit run dashboard.py --server.address 0.0.0.0 --server.port 8501 --server.headless true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
